# docker/cron.daily
#!/bin/sh
cd /var/www/app;

php bin/console slideo:clear:downloads;